# presentations/systemd-2016/demo/demo.service
[Unit]
Description=Demo unit for systemd

[Service]
Type=oneshot
ExecStart=/usr/bin/echo "Systemd unit activated. Hello Hackeriet."
