# presentations/systemd-2016/demo/demo-stress.service
[Unit]
Description=Stress test CPU

[Service]
Slice=demo.slice
ExecStart=/usr/bin/stress -c 5
