# presentations/systemd-2016/demo/demo-path.path
[Unit]
Description=Demonstrate systemd path units

[Path]
DirectoryNotEmpty=/tmp/hackeriet
Unit=demo.service
