# demo/demo-limits.slice
[Unit]
Description=Limited resources demo
DefaultDependencies=no
Before=slices.target

[Slice]
CPUQuota=10%
