# demo/demo-error@.service
[Unit]
Description=Demonstrate failing units
OnFailure=demo-notify@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/false
