# demo/demo-timer.service
[Unit]
Description=Demo unit for systemd timers

[Service]
Type=oneshot
ExecStart=/usr/bin/echo "Systemd timer fired! Hello Hackeriet."
